# .config/yadm/alt/.config/archinstall##d.arch/etc/profile.d/40-path.sh
#!/usr/bin/env sh

export PATH=$HOME/.local/bin:$PATH
export PATH=$CARGO_HOME/bin:$PATH
export PATH=$GOPATH/bin:$PATH
