# .config/yadm/root/opt##o.Darwin/local/etc/profile.d/70-work.sh
export QAPM_HOME="$HOME/Developer/QAPM"
export INTRANET_PROXY_JUMP_SERVER='mixed-main-22'
export TKE_1_8_PROXY_JUMP_SERVER='tke1.8'