# .config/yadm/alt/.config/install##o.Darwin/opt/local/etc/profile.d/60-macport.sh
export C_INCLUDE_PATH=/opt/local/include
export CPLUS_INCLUDE_PATH=/opt/local/include
export LIBRARY_PATH=/opt/local/lib
export JAVA_HOME='/Library/Java/JavaVirtualMachines/openjdk8-zulu/Contents/Home/'
