# Library/LaunchAgents/custom.program.env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>custom.program.env</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>
        /bin/launchctl setenv GOPROXY "https://goproxy.cn,direct"
        /bin/launchctl setenv RUSTUP_DIST_SERVER "https://rsproxy.cn"
        /bin/launchctl setenv RUSTUP_UPDATE_ROOT "https://rsproxy.cn/rustup"
        /bin/launchctl setenv RUSTC_WRAPPER "sccache"
      </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
  </dict>
</plist>