# Library/LaunchAgents/custom.system.env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>custom.system.env</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>
        /bin/launchctl setenv EDITOR "helix";
        /bin/launchctl setenv PAGER "less";
      </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!--
    <key>StandardErrorPath</key>
    <string>/tmp/setting-env.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/setting-env.log</string>
    <key>Debug</key>
    <true/>
    -->
  </dict>
</plist>