# Library/LaunchAgents/custom.program.env.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>custom.program.env</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>
        /bin/launchctl setenv GOPROXY "https://goproxy.cn,direct";
        /bin/launchctl setenv RUSTUP_DIST_SERVER "https://rsproxy.cn";
        /bin/launchctl setenv RUSTUP_UPDATE_ROOT "https://rsproxy.cn/rustup";
        /bin/launchctl setenv RUSTC_WRAPPER "sccache";
        /bin/launchctl setenv HOMEBREW_NO_ANALYTICS "1";
{{- if .gpgagent }}
        if [ "${gnupg_SSH_AUTH_SOCK_by:-0}" -ne $$ ]; then
          /bin/launchctl setenv SSH_AUTH_SOCK "$(gpgconf --list-dirs agent-ssh-socket)";
        fi
        /bin/launchctl unsetenv SSH_AGENT_PID;
{{- end }}
      </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!--
    <key>StandardErrorPath</key>
    <string>/tmp/setting-env.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/setting-env.log</string>
    <key>Debug</key>
    <true/>
    -->
  </dict>
</plist>